# backup-volumes.timer
[Unit]
Description=Backup podman volumes

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
