# backup-volumes.service
[Unit]
Description=Backup podman volumes

[Service]
Type=oneshot
ExecStart=%h/bin/backup-volumes.sh -o %h/backups
